# root/etc/systemd/system/rh-php73-php-fpm.service.d/nethserver-rh-php73-php-fpm.conf
[Service]
RuntimeDirectory=rh-php73-php-fpm
RuntimeDirectoryMode=0755
